# templates/netbox-rqworker.conf.j2
# {{ ansible_managed }}

[program:netbox-rqworker]
command = {{ netbox_user_home_directory }}/rqworker_process.sh
user = {{ netbox_user }}
