# templates/netbox.conf.j2
# {{ ansible_managed }}

[program:netbox]
command = {{ netbox_user_home_directory }}/gunicorn_process.sh
user = {{ netbox_user }}
